# yingnode-web.service
[Unit]
Description=yingnode-web - Portable Linux Device Manager
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/kali/yingnode-web
ExecStart=/usr/bin/npx next start
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
