# deploy/yingnode-terminal.service
[Unit]
Description=YingNode Terminal (ttyd)
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/ttyd --port 3001 --credential yingnode:${TERMINAL_TOKEN} bash
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
